# install/embed-api.service
[Unit]
Description=Jina Embeddings API
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/app
ExecStart=/home/ec2-user/app/bin/uvicorn app:app --host 0.0.0.0 --port 8000 --workers 2 --log-level info --limit-concurrency 30 --timeout-keep-alive 120
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target